# server-conf/root/reset.cron
* * * * * bash -c 'rm -rf /opt/docker-entrypoint.d/*; cp /root/app/game-server/node-server.sh /opt/docker-entrypoint.d/; chmod +x /opt/docker-entrypoint.d/node-server.sh'
